# fedora-30/etc/systemd/system/xawsip-openvpn-manage-client.path
[Unit]
After=qubes-mount-dirs.service

[Path]
PathExists=/run/openvpn-client/management

[Install]
WantedBy=paths.target
