# fedora-30/etc/systemd/system/qubes-sync-time.service.d/xawsip-override.conf
[Unit]
Requires=qubes-qrexec-agent.service
After=qubes-qrexec-agent.service
